# docker/deploy/sensor_emulator_task.conf
[program:sensor_emulator]
directory=/usr/src/app/
command=python bin/start_emulation.py
redirect_stderr=true
