# docker/deploy/sensor_emulator_task.conf
[program:sensor_emulator]
directory=/usr/src/app/
command=python sensorizer/bin/emulate_full_plants.py
redirect_stderr=true